# slideshow@.service
[Unit]
Description=Photo Slideshow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/%i/slideshow.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
